# etc/udev/rules.d/70-snap.woe-usb.rules
# This file is automatically generated.
# optical-drive
KERNEL=="scd[0-9]*", TAG+="snap_woe-usb_woeusb"
# optical-drive
KERNEL=="sr[0-9]*", TAG+="snap_woe-usb_woeusb"
# optical-drive
SUBSYSTEM=="scsi_generic", SUBSYSTEMS=="scsi", ATTRS{type}=="4|5", TAG+="snap_woe-usb_woeusb"
# raw-usb
SUBSYSTEM=="tty", ENV{ID_BUS}=="usb", TAG+="snap_woe-usb_woeusb"
# raw-usb
SUBSYSTEM=="usb", TAG+="snap_woe-usb_woeusb"
# raw-usb
SUBSYSTEM=="usbmisc", TAG+="snap_woe-usb_woeusb"
TAG=="snap_woe-usb_woeusb", SUBSYSTEM!="module", SUBSYSTEM!="subsystem", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_woe-usb_woeusb $devpath $major:$minor"
# optical-drive
KERNEL=="scd[0-9]*", TAG+="snap_woe-usb_woeusbgui"
# optical-drive
KERNEL=="sr[0-9]*", TAG+="snap_woe-usb_woeusbgui"
# optical-drive
SUBSYSTEM=="scsi_generic", SUBSYSTEMS=="scsi", ATTRS{type}=="4|5", TAG+="snap_woe-usb_woeusbgui"
# raw-usb
SUBSYSTEM=="tty", ENV{ID_BUS}=="usb", TAG+="snap_woe-usb_woeusbgui"
# raw-usb
SUBSYSTEM=="usb", TAG+="snap_woe-usb_woeusbgui"
# raw-usb
SUBSYSTEM=="usbmisc", TAG+="snap_woe-usb_woeusbgui"
TAG=="snap_woe-usb_woeusbgui", SUBSYSTEM!="module", SUBSYSTEM!="subsystem", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_woe-usb_woeusbgui $devpath $major:$minor"
